# cross/icinga2/PLIST
etc/bash_completion.d/icinga2
etc/icinga2/conf.d/app.conf
etc/icinga2/conf.d/commands.conf
etc/icinga2/conf.d/downtimes.conf
etc/icinga2/conf.d/groups.conf
etc/icinga2/conf.d/hosts.conf
etc/icinga2/conf.d/notifications.conf
etc/icinga2/conf.d/services.conf
etc/icinga2/conf.d/templates.conf
etc/icinga2/conf.d/timeperiods.conf
etc/icinga2/conf.d/users.conf
etc/icinga2/constants.conf
etc/icinga2/features-available/api.conf
etc/icinga2/features-available/checker.conf
etc/icinga2/features-available/command.conf
etc/icinga2/features-available/compatlog.conf
etc/icinga2/features-available/debuglog.conf
etc/icinga2/features-available/elasticsearch.conf
etc/icinga2/features-available/gelf.conf
etc/icinga2/features-available/graphite.conf
etc/icinga2/features-available/icingadb.conf
etc/icinga2/features-available/ido-mysql.conf
etc/icinga2/features-available/ido-pgsql.conf
etc/icinga2/features-available/influxdb.conf
etc/icinga2/features-available/influxdb2.conf
etc/icinga2/features-available/livestatus.conf
etc/icinga2/features-available/mainlog.conf
etc/icinga2/features-available/notification.conf
etc/icinga2/features-available/opentsdb.conf
etc/icinga2/features-available/perfdata.conf
etc/icinga2/features-available/syslog.conf
etc/icinga2/features-enabled/checker.conf
etc/icinga2/features-enabled/mainlog.conf
etc/icinga2/features-enabled/notification.conf
etc/icinga2/icinga2.conf
etc/icinga2/scripts/mail-host-notification.sh
etc/icinga2/scripts/mail-service-notification.sh
etc/icinga2/zones.conf
etc/icinga2/zones.d/README
etc/init.d/icinga2
etc/logrotate.d/icinga2
etc/sysconfig/icinga2
lib/icinga2/libmysql_shim.so
lib/icinga2/libmysql_shim.so.2.15.1
lib/icinga2/libpgsql_shim.so
lib/icinga2/libpgsql_shim.so.2.15.1
lib/icinga2/prepare-dirs
lib/icinga2/safe-reload
lib/icinga2/sbin/icinga2
sbin/icinga2
share/doc/icinga2/AUTHORS
share/doc/icinga2/CHANGELOG.md
share/doc/icinga2/COPYING
share/doc/icinga2/NEWS
share/doc/icinga2/README.md
share/doc/icinga2/markdown/01-about.md
share/doc/icinga2/markdown/02-installation.md
share/doc/icinga2/markdown/03-monitoring-basics.md
share/doc/icinga2/markdown/04-configuration.md
share/doc/icinga2/markdown/05-service-monitoring.md
share/doc/icinga2/markdown/06-distributed-monitoring.md
share/doc/icinga2/markdown/07-agent-based-monitoring.md
share/doc/icinga2/markdown/08-advanced-topics.md
share/doc/icinga2/markdown/09-object-types.md
share/doc/icinga2/markdown/10-icinga-template-library.md
share/doc/icinga2/markdown/11-cli-commands.md
share/doc/icinga2/markdown/12-icinga2-api.md
share/doc/icinga2/markdown/13-addons.md
share/doc/icinga2/markdown/14-features.md
share/doc/icinga2/markdown/15-troubleshooting.md
share/doc/icinga2/markdown/16-upgrading-icinga-2.md
share/doc/icinga2/markdown/17-language-reference.md
share/doc/icinga2/markdown/18-library-reference.md
share/doc/icinga2/markdown/19-technical-concepts.md
share/doc/icinga2/markdown/20-script-debugger.md
share/doc/icinga2/markdown/21-development.md
share/doc/icinga2/markdown/22-selinux.md
share/doc/icinga2/markdown/23-migrating-from-icinga-1x.md
share/doc/icinga2/markdown/24-appendix.md
share/doc/icinga2/markdown/images/addons/dashing_icinga2.png
share/doc/icinga2/markdown/images/addons/icinga_certificate_monitoring.png
share/doc/icinga2/markdown/images/addons/icinga_reporting.png
share/doc/icinga2/markdown/images/addons/icingaweb2_businessprocess.png
share/doc/icinga2/markdown/images/addons/icingaweb2_grafana.png
share/doc/icinga2/markdown/images/addons/icingaweb2_graphite.png
share/doc/icinga2/markdown/images/addons/icingaweb2_maps.png
share/doc/icinga2/markdown/images/addons/nano-syntax.png
share/doc/icinga2/markdown/images/addons/vim-syntax.png
share/doc/icinga2/markdown/images/advanced-topics/flapping-state-graph.png
share/doc/icinga2/markdown/images/advanced-topics/icinga2_external_checks_freshness_icingaweb2.png
share/doc/icinga2/markdown/images/advanced-topics/icingaweb2_downtime_handled.png
share/doc/icinga2/markdown/images/api/icinga2_api_powershell_ise.png
share/doc/icinga2/markdown/images/configuration/icinga_web_local_server.png
share/doc/icinga2/markdown/images/development/windows_boost_build_dev_cmd.png
share/doc/icinga2/markdown/images/development/windows_builds_gitlab_pipeline.png
share/doc/icinga2/markdown/images/development/windows_powershell_posh_git.png
share/doc/icinga2/markdown/images/development/windows_visual_studio_installer_01.png
share/doc/icinga2/markdown/images/development/windows_visual_studio_installer_02.png
share/doc/icinga2/markdown/images/development/windows_visual_studio_installer_03.png
share/doc/icinga2/markdown/images/development/windows_visual_studio_tabs_c++.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_agent_checks_command_endpoint.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_endpoints.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_roles.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_satellite_config_sync.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_scenario_ha_masters_with_agents.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_scenarios_master_satellites_agents.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_scenarios_master_with_agents.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_monitoring_zones.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_windows_client_disk_icingaweb2.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_windows_nscp_api_drivesize_icingaweb2.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_distributed_windows_nscp_counter_icingaweb2.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_running_service.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_installer_01.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_installer_02.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_installer_03.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_installer_04.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_installer_05.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_01.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_02.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_02_global_zone.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_03.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_04.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_06_finish_no_ticket.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_06_finish_with_ticket.png
share/doc/icinga2/markdown/images/distributed-monitoring/icinga2_windows_setup_wizard_examine_config.png
share/doc/icinga2/markdown/images/icingadb/icingadb-architecture.png
share/doc/icinga2/markdown/images/icingadb/icingadb-daemon.png
share/doc/icinga2/markdown/images/icingadb/icingadb-icinga2.png
share/doc/icinga2/markdown/images/icingadb/icingadb-redis.png
share/icinga2-ido-mysql/schema/mysql.sql
share/icinga2-ido-mysql/schema/upgrade/2.0.2.sql
share/icinga2-ido-mysql/schema/upgrade/2.1.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.11.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.12.7.sql
share/icinga2-ido-mysql/schema/upgrade/2.13.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.13.3.sql
share/icinga2-ido-mysql/schema/upgrade/2.2.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.3.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.4.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.5.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.6.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.8.0.sql
share/icinga2-ido-mysql/schema/upgrade/2.8.1.sql
share/icinga2-ido-pgsql/schema/pgsql.sql
share/icinga2-ido-pgsql/schema/upgrade/2.0.2.sql
share/icinga2-ido-pgsql/schema/upgrade/2.1.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.2.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.3.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.4.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.5.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.6.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.8.0.sql
share/icinga2-ido-pgsql/schema/upgrade/2.8.1.sql
share/icinga2/include/command-icinga.conf
share/icinga2/include/command-nscp-local.conf
share/icinga2/include/command-plugins-manubulon.conf
share/icinga2/include/command-plugins-windows.conf
share/icinga2/include/command-plugins.conf
share/icinga2/include/hangman
share/icinga2/include/itl
share/icinga2/include/manubulon
share/icinga2/include/nscp
share/icinga2/include/plugins
share/icinga2/include/plugins-contrib
share/icinga2/include/plugins-contrib.d/big-data.conf
share/icinga2/include/plugins-contrib.d/databases.conf
share/icinga2/include/plugins-contrib.d/hardware.conf
share/icinga2/include/plugins-contrib.d/icingacli.conf
share/icinga2/include/plugins-contrib.d/ipmi.conf
share/icinga2/include/plugins-contrib.d/logmanagement.conf
share/icinga2/include/plugins-contrib.d/metrics.conf
share/icinga2/include/plugins-contrib.d/network-components.conf
share/icinga2/include/plugins-contrib.d/network-services.conf
share/icinga2/include/plugins-contrib.d/operating-system.conf
share/icinga2/include/plugins-contrib.d/raid-controller.conf
share/icinga2/include/plugins-contrib.d/smart-attributes.conf
share/icinga2/include/plugins-contrib.d/storage.conf
share/icinga2/include/plugins-contrib.d/systemd.conf
share/icinga2/include/plugins-contrib.d/virtualization.conf
share/icinga2/include/plugins-contrib.d/vmware.conf
share/icinga2/include/plugins-contrib.d/web.conf
share/icinga2/include/windows-plugins
share/man/man8/icinga2.8
